# Systemd/API.service
[Unit]
Description=Programa API

[Service]
ExecStart=/ruta/al/programa1
Restart=always

[Install]
WantedBy=multi-user.target